# dji-h264-stream.service
[Unit]
Description=DJI Osmo Pocket 3 H264 4K 50fps Stream to HDMI with Audio
After=network.target sound.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/dji-stream.sh
Restart=always
RestartSec=3
Nice=-10
IOSchedulingClass=realtime

[Install]
WantedBy=multi-user.target
