# hdmi-button.service
[Unit]
Description=HDMI Mode Toggle Button Daemon
After=multi-user.target
Documentation=https://github.com/stephanebhiri/DJI_OSMOPOCKET3_TO_HDMI_4K_60P_50P

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/gpio-button-daemon.py 100
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
